# docker/php-fpm/crontab/pms-pi
PATH=/usr/local/bin
BIN_CONSOLE_PATH=/application/bin/console

# At 01:00 on every 3rd day-of-week.
0 1 * * */3 www-data php "$BIN_CONSOLE_PATH" cron:make-config-backup "/backup/pms" config

# At 02:00 on every 3rd day-of-week.
0 2 * * */3 www-data php "$BIN_CONSOLE_PATH" cron:make-backup "/backup/pms" database files
